# deploy/systemd/user/ai-memory-rotate.timer
[Unit]
Description=Run AI Memory Rotation Weekly With Catch-Up

[Timer]
OnBootSec=5min
OnCalendar=weekly
Persistent=true
Unit=ai-memory-rotate.service

[Install]
WantedBy=timers.target
